# maconfig/file-config/rc.local
#!/bin/bash
echo "never" | sudo tee /sys/kernel/mm/transparent_hugepage/enabled  
echo "never" | sudo tee /sys/kernel/mm/transparent_hugepage/defrag
sysctl vm.swappiness=10
exit 0
